# deploy/slayseed.service
[Unit]
Description=SlaySeed Spring Boot Application
After=network.target mysql.service redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/slayseed
ExecStart=/usr/bin/java -jar /opt/slayseed/backend/slayseed.jar
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/slayseed/backend.log
StandardError=append:/var/log/slayseed/backend-error.log

# 环境变量 - 按需修改
Environment="JAVA_OPTS=-Xms512m -Xmx1024m
Environment="JASYPT_ENCRYPTOR_KEY=your-encryptor-key"
Environment="SPRING_DATASOURCE_URL=jdbc:mysql://127.0.0.1:3306/slayseed?useUnicode=true&characterEncoding=utf8&useSSL=false&serverTimezone=Asia/Shanghai&allowPublicKeyRetrieval=true"
Environment="SPRING_DATASOURCE_USERNAME=root"
Environment="SPRING_DATASOURCE_PASSWORD=your-db-password"
Environment="SPRING_DATA_REDIS_HOST=127.0.0.1"
Environment="SPRING_DATA_REDIS_PORT=6379"
Environment="SPRING_DATA_REDIS_PASSWORD="
Environment="SERVER_PORT=8080"
Environment="CORS_ALLOWED_ORIGINS=*"

[Install]
WantedBy=multi-user.target